# deploy/onvifscanner.service
[Unit]
Description=ONVIFscanner (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=onvifscanner
Group=onvifscanner
WorkingDirectory=/opt/onvifscanner
EnvironmentFile=/etc/onvifscanner/onvifscanner.env
ExecStart=/usr/local/bin/onvifscanner-start
Restart=on-failure
RestartSec=2
WatchdogSec=30
TimeoutStartSec=60
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/onvifscanner

[Install]
WantedBy=multi-user.target
